# app/config/supervisord/supervisord.conf
[supervisord]
logfile=/dev/stdout
pidfile=/tmp/supervisord.pid

[program:socat]
command=/app/bin/wrap-socat.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout

[program:newrelic-daemon]
command=/app/bin/wrap-newrelic-daemon.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
depends_on=socat
